# systemd/forgewatch-indicator.service
# Reference copy — do not install manually.
# Use `forgewatch service install` which resolves the correct
# ExecStart path for your installation.

[Unit]
Description=ForgeWatch — Panel Indicator
After=forgewatch.service
Wants=forgewatch.service

[Service]
Type=simple
ExecStart=%h/.local/bin/forgewatch-indicator
Restart=on-failure
RestartSec=10

# Security hardening (lighter than daemon — no file writes needed)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=default.target
