# systemd/forgewatch.service
# Reference copy — do not install manually.
# Use `forgewatch service install` which resolves the correct
# ExecStart path for your installation.

[Unit]
Description=ForgeWatch — GitHub PR Monitor
After=network-online.target dbus.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=%h/.local/bin/forgewatch
ExecReload=kill -HUP $MAINPID
Restart=on-failure
RestartSec=10
Environment=GITHUB_TOKEN=

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.config/forgewatch

[Install]
WantedBy=default.target
