# forgewatch/cli/systemd/forgewatch-indicator.service
[Unit]
Description=ForgeWatch — Panel Indicator
After=forgewatch.service
Wants=forgewatch.service

[Service]
Type=simple
ExecStart=@@FORGEWATCH_INDICATOR_EXEC@@
Restart=on-failure
RestartSec=10

# Security hardening (lighter than daemon — no file writes needed)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=default.target
